# automation/deploy/com.tim-bankrupt.igpublish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent: poll the publish queue every 60s and auto-publish due posts.
     Why 60s (not 15min): publish_at carries a human-ish jittered minute (e.g. 15:02); a coarse poll
     would round the ACTUAL post time up to the next tick and waste the jitter. A 60s poll fires at
     the planned minute. It's free: when nothing is due the run is a sub-second local file read with
     NO Graph API call. StartInterval is offset from load time, so ticks are never on a round :00.
     Install: edit the two absolute paths below to your checkout, then:
       cp automation/deploy/com.tim-bankrupt.igpublish.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.tim-bankrupt.igpublish.plist
     A run on wake catches any post that came due while the Mac was asleep (publishes a bit late,
     never silently dropped). See automation/deploy/README.md. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.tim-bankrupt.igpublish</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/tatsumiyuuki/Desktop/dev/midjourney-explore/automation/deploy/run_publish.sh</string>
  </array>
  <key>StartInterval</key><integer>60</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/tatsumiyuuki/Desktop/dev/midjourney-explore/automation/deploy/igpub.log</string>
  <key>StandardErrorPath</key><string>/Users/tatsumiyuuki/Desktop/dev/midjourney-explore/automation/deploy/igpub.log</string>
</dict>
</plist>
